# engine/supervisord.conf
; Standalone supervisord config for the mt5-wine container.
; Programs are declared here; /etc/supervisor/supervisord.conf (Debian default)
; includes /etc/supervisor/conf.d/*.conf, so this file is picked up automatically.

[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ── Virtual framebuffer ───────────────────────────────────────────────────────
[program:xvfb]
command=Xvfb :99 -screen 0 1024x768x16 -nolisten tcp -nolisten inet6
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Installer display (shared :99) → noVNC websocket ────────────────────────
; Exposes the bootstrap desktop so the frontend can render it during install.
[program:installer-x11vnc]
command=x11vnc -display :99 -rfbport 5999 -nopw -forever -shared -noxdamage -q -o /mt5-fleet/config/installer-x11vnc.log
autorestart=true
priority=11
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:installer-websockify]
command=websockify 0.0.0.0:6799 localhost:5999
autorestart=true
priority=12
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Runtime installer (local builds) ─────────────────────────────────────────
; For local compose we skip install_reference.sh at docker-build time so the UI
; can start immediately. This job runs once at container startup, logs to both
; container stdout and /mt5-fleet/config/install.log, and exits after success.
[program:installer]
command=/bin/sh -lc 'mkdir -p /mt5-fleet/config; if [ ! -f /mt5-fleet/.installed ]; then : > /mt5-fleet/config/install.log; /opt/mt5/scripts/install_reference.sh 2>&1 | tee -a /mt5-fleet/config/install.log; else echo "[install] already installed"; fi'
autorestart=false
priority=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Control API (manages workers, called by the api service) ─────────────────
[program:control-api]
command=/opt/mt5/control_api/server
autorestart=true
priority=20
startsecs=2
environment=FLEET_DIR="/mt5-fleet",WINEPREFIX="/mt5-fleet/wineprefix",DISPLAY=":99",WINEDEBUG="-all",PORT="18810"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
