# engine/supervisord.conf
; Standalone supervisord config for the mt5-wine container.
; Programs are declared here; /etc/supervisor/supervisord.conf (Debian default)
; includes /etc/supervisor/conf.d/*.conf, so this file is picked up automatically.

[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ── Installer virtual display (Xvnc on :99) ──────────────────────────────────
; Replaced Xvfb with Xvnc so the installer desktop is directly VNC-accessible.
[program:installer-xvnc]
command=/bin/sh -c 'rm -f /tmp/.X99-lock /tmp/.X11-unix/X99; exec Xvnc :99 -geometry 1920x1080 -depth 24 -rfbport 5999 -SecurityTypes None -AcceptSetDesktopSize -AlwaysShared -nolisten tcp -nolisten inet6'
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Installer window manager (autotileWM) ─────────────────────────────────────
[program:installer-wm]
command=autotilewm
autorestart=true
priority=11
startsecs=1
environment=DISPLAY=":99",XDG_RUNTIME_DIR="/tmp"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Installer websockify (WebSocket → VNC proxy) ─────────────────────────────
[program:installer-websockify]
command=websockify 0.0.0.0:6799 localhost:5999
autorestart=true
priority=12
startsecs=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Control API (manages workers, called by the api service) ─────────────────
[program:control-api]
command=/opt/mt5/control_api/server
autorestart=true
priority=18
startsecs=2
environment=FLEET_DIR="/mt5-fleet",WINEPREFIX="/mt5-fleet/wineprefix",DISPLAY=":99",WINEDEBUG="-all",PORT="18810"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
