# engine/supervisord.conf
; Standalone supervisord config for the mt5-wine container.

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=20MB
logfile_backups=3
loglevel=warn
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ── Reference virtual display (Xvnc on :99) ──────────────────────────────────
; The reference display is always live. The user interacts with MT5 here to
; configure brokers, install EAs, etc. New terminals clone from this reference.
[program:reference-xvnc]
command=/bin/sh -c 'rm -f /tmp/.X99-lock /tmp/.X11-unix/X99; exec Xvnc :99 -geometry 1920x1080 -depth 24 -rfbport 5999 -SecurityTypes None -AcceptSetDesktopSize -AlwaysShared -nolisten tcp -nolisten inet6'
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/reference-xvnc.out.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/reference-xvnc.err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3

; ── Reference window manager (AwesomeWM on :99) ─────────────────────────────
[program:reference-wm]
command=/bin/sh -c 'exec awesome -c /opt/mt5/scripts/awesome.rc.lua'
autorestart=true
priority=11
startsecs=2
environment=DISPLAY=":99",XDG_RUNTIME_DIR="/tmp",HOME="/root"
stdout_logfile=/var/log/supervisor/reference-wm.out.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/reference-wm.err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3

; ── Reference websockify (WebSocket → VNC proxy for :99) ─────────────────────
[program:reference-websockify]
command=websockify 0.0.0.0:6799 localhost:5999
autorestart=true
priority=12
startsecs=1
stdout_logfile=/var/log/supervisor/reference-websockify.out.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/reference-websockify.err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3

; ── Reference installer (runs once on first boot) ─────────────────────────────
; install_reference.sh is idempotent: exits immediately if already installed.
; autorestart=false so it does not re-run after a successful exit.
[program:reference-install]
command=/opt/mt5/scripts/install_reference.sh
autorestart=false
priority=13
startsecs=0
exitcodes=0
environment=FLEET_DIR="/my5fleet",WINEPREFIX="/my5fleet/wineprefix",DISPLAY=":99",WINEDEBUG="-all"
stdout_logfile=/var/log/supervisor/reference-install.out.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/reference-install.err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3

; ── Reference MT5 (runs on :99 after install completes) ──────────────────────
; reference_mt5.sh waits for terminal64.exe to appear then exec's MT5.
; autorestart=true so MT5 restarts if it crashes.
[program:reference-mt5]
command=/opt/mt5/scripts/reference_mt5.sh
autorestart=true
priority=20
startsecs=10
environment=FLEET_DIR="/my5fleet",WINEPREFIX="/my5fleet/wineprefix",DISPLAY=":99",WINEDEBUG="-all"
stdout_logfile=/var/log/supervisor/reference-mt5.out.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/reference-mt5.err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3

; ── Control API (manages workers, called by the api service) ─────────────────
[program:control-api]
command=/opt/mt5/control_api/server
autorestart=true
priority=18
startsecs=2
environment=FLEET_DIR="/my5fleet",WINEPREFIX="/my5fleet/wineprefix",DISPLAY=":99",WINEDEBUG="-all",PORT="18810"
stdout_logfile=/var/log/supervisor/control-api.out.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/control-api.err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=3